# deploy/systemd/tracker.service
[Unit]
Description=Amazon Tracker v2
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/amazon-tracker
EnvironmentFile=/opt/amazon-tracker/.env
ExecStart=/opt/amazon-tracker/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 1
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
